# init/sletatru-search.conf.example
# /etc/init/sletatru-search.conf

start on started network
stop on stopping network

respawn
expect daemon

exec sletatru-search